# systemd-ai.service
[Unit]
Description=Systemd_AI
After=network.target NetworkManager.service

[Service]
Type=simple
ExecStart=/usr/share/Systemd_AI/python_env/bin/python3 /usr/share/Systemd_AI/systemd_analyzer.py
User=root
Group=root
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
