# systemd-ai.timer
[Unit]
Description=systemd analyzer 12h timer
[Timer]
OnBootSec=1min
OnUnitActiveSec=12h
Unit=systemd-ai.service

[Install]
WantedBy=timers.target
